# apps/vaultwarden/vaultwarden.container
[Unit]
Description=Vaultwarden server

[Service]
ExecStartPre=sh vaultwarden/vaultwarden.bootstrap

[Container]
Image=docker.io/vaultwarden/server
ContainerName=vaultwarden

Environment=SIGNUPS_ALLOWED=false
Environment=INVITATIONS_ALLOWED=false
Environment=ROCKET_PORT=8080

Volume=${VOLUME_DIR}/vaultwarden/data:/data
