# server/caddy.container
[Unit]
Description=Caddy Reverse Proxy

[Service]
ExecStartPre=sh server/caddy.bootstrap

[Container]
Image=docker.io/cilenco/caddy-ionos-dns

HostName=caddy
ContainerName=caddy
AutoUpdate=registry
Network=caddy

Environment=DOMAIN=${DOMAIN}

Secret=IONOS_PREFIX,type=env,target=IONOS_PREFIX
Secret=IONOS_SECRET,type=env,target=IONOS_SECRET

PublishPort=80:8080
PublishPort=8883:8883
PublishPort=443:4430

Volume=./Caddyfile:/etc/caddy/Caddyfile:ro

Volume=${VOLUMES_DIR}/caddy/certs:/certs:idmap
Volume=${VOLUMES_DIR}/caddy/data:/data:idmap
