# ansible/roles/autorestic/autorestic.timer
[Unit]
Description=Autorestic backup timer

[Timer]
Unit=autorestic.service
OnCalendar=*:0/15

[Install]
WantedBy=default.target
